# docker-builds/01_base/etc/cron.d/anacron
# /etc/cron.d/anacron: crontab entries for the anacron package
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
LANGUAGE=en_US.UTF-8
LANG=en_US.UTF-8
LC_ALL=en_US.UTF-8
30 * * * *   root . /etc/profile; /etc/init.d/anacron start
