# docker-builds/02_ldap/etc/cron.d/backup_ldap
0 4 * * *   root /usr/local/bin/backup-ldap 2>&1 >>/var/log/backup-ldap.log
