# docker-builds/02_ldap/etc/cron.d/backup_ldap
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
LANGUAGE=en_US.UTF-8
LANG=en_US.UTF-8
LC_ALL=en_US.UTF-8
0 4 * * *   root . /etc/profile; backup-ldap 2>&1 >>/var/log/backup-ldap.log
